# detections/endpoint/crontab_job_scheduling_(unix)-splunk-edr.yml
id: '1132.24776'
title: Crontab Job Scheduling (Unix)
description: 'Usecase which would detect when a cronjob is created using crontab on
  Unix hosts. Atomics T1053.003 Test#1 Atomics T1053.003 Test#2 Atomics T1053.003
  Test#3 Threat Actor Association: Sandworm (UAC-0165), WatchDog'
logic_format: Splunk
logic: '`get_endpoint_data` `get_endpoint_data_edr` ("crontab" "-e") OR TERM(echo)
  OR TERM(>) OR TERM(>>) | regex "(((\*|\*\/\d+|\d+)\s){3,})|\>\>?.*?(etc.cron|var.spool.cron)|crontab.*?\-e|echo.+?/cron"
  | table _time, host, user process, process_* | bin span=1s | stats values(*) as
  * by _time, host '
techniques:
- execution:scheduled task/job:cron
- persistence:scheduled task/job:cron
- privilege-escalation:scheduled task/job:cron
technique_id:
- T1053.003
data_category:
- Process command-line parameters
- EDR Logs
references:
- https://materials.rangeforce.com/tutorial/2020/04/17/Cron-Privilege-Escalation/
