# detections/endpoint/crontab_job_scheduling_(unix)-snowflake-crowdstrikefdr_process.yml
id: '1132.14563'
title: Crontab Job Scheduling (Unix)
description: 'Usecase which would detect when a cronjob is created using crontab on
  Unix hosts. Atomics T1053.003 Test#1 Atomics T1053.003 Test#2 Atomics T1053.003
  Test#3 Threat Actor Association: Sandworm (UAC-0165), WatchDog'
logic_format: snowflake
logic: 'select * from crowdstrikefdr_process where event_time > dateadd(hour, -2,
  sysdate()) and (event_platform = ''Lin'' or event_platform = ''Mac'') and regexp_like(process,
  ''.*((((\\\*|\\\*\\\/\\\d+|\\\d+)\\\s){3,})|\\\>\\\>?.*?(etc.cron|var.spool.cron)|crontab.*?\\\-e|echo.+?/cron).*'',
  ''i'') and (process_path ilike ''%echo%'' or process_path ilike ''%crontab%'') '
techniques:
- execution:scheduled task/job:cron
- persistence:scheduled task/job:cron
- privilege-escalation:scheduled task/job:cron
technique_id:
- T1053.003
data_category:
- EDR Logs
references:
- https://materials.rangeforce.com/tutorial/2020/04/17/Cron-Privilege-Escalation/
